# infra/server/systemd/livraison-backup.timer
[Unit]
Description=Run Livraison database backup nightly

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
