# infra/server/systemd/livraison-backup.service
[Unit]
Description=Livraison database backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/livraison/infra/server/scripts/backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
