# infra/server/systemd/livraison-health.service
[Unit]
Description=Livraison health check and self-heal
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/livraison/infra/server/scripts/healthcheck.sh
