# infra/server/systemd/livraison-monitor.service
[Unit]
Description=Livraison host resource monitor (disk, memory, containers)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/livraison/infra/server/scripts/monitor.sh
